# mcp-servers/agents/systemd/content-agent.service
[Unit]
Description=AI CoS Content Agent
After=network.target sync-agent.service
Wants=sync-agent.service

[Service]
Type=simple
WorkingDirectory=/opt/agents
EnvironmentFile=/opt/agents/.env
ExecStartPre=/bin/bash /opt/agents/scripts/preflight.sh content-agent
ExecStart=/opt/agents/.venv/bin/python3 -m content.server
Restart=always
RestartSec=5
MemoryMax=768M

[Install]
WantedBy=multi-user.target
