# mcp-servers/agents/systemd/sync-agent.service
[Unit]
Description=AI CoS Sync Agent (Gateway)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/agents
EnvironmentFile=/opt/agents/.env
ExecStartPre=/bin/bash /opt/agents/scripts/preflight.sh sync-agent
ExecStart=/opt/agents/.venv/bin/python3 -m sync.server
Restart=always
RestartSec=5
MemoryMax=512M

[Install]
WantedBy=multi-user.target
