# mcp-servers/agents/systemd/sync-agent.service
[Unit]
Description=AI CoS Sync Agent (Gateway)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/agents
EnvironmentFile=-/opt/agents/.env
ExecStart=/opt/agents/.venv/bin/python3 -m sync.server
Restart=always
RestartSec=5
WatchdogSec=120
MemoryMax=512M
StandardOutput=append:/opt/agents/logs/sync.log
StandardError=append:/opt/agents/logs/sync.log

[Install]
WantedBy=multi-user.target
