# mcp-servers/agents/systemd/web-agent.service
[Unit]
Description=AI CoS Web Agent
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/agents
EnvironmentFile=/opt/agents/.env
ExecStartPre=/bin/bash /opt/agents/scripts/preflight.sh web-agent
ExecStart=/opt/agents/.venv/bin/python3 -m web.server
Restart=always
RestartSec=5
MemoryMax=1G
StandardOutput=append:/opt/agents/logs/web.log
StandardError=append:/opt/agents/logs/web.log

[Install]
WantedBy=multi-user.target
